# deb/install/postinst
#!/bin/bash

case "${1}" in
  configure)
    echo "Enabling waggle-bk-registration service"
    systemctl enable waggle-bk-registration.service
    ;;
esac
