# deb/install/prerm
#!/bin/bash

systemctl stop waggle-bk-registration.service || true
systemctl disable waggle-bk-registration.service || true
