# deb/reg/postinst
#!/bin/bash

case "${1}" in
  configure)
    echo "Enabling & starting waggle-registration service"
    systemctl enable waggle-registration.service
    ;;
esac
