# deb/reg/prerm
#!/bin/bash

systemctl stop sagebk-registration.service || true
systemctl disable sagebk-registration.service || true
