# deb/reg/prerm
#!/bin/bash

systemctl stop waggle-registration.service || true
systemctl disable waggle-registration.service || true
