# deb/reverse/prerm
#!/bin/bash

systemctl stop sagebk-reverse-tunnel.service || true
systemctl disable sagebk-reverse-tunnel.service || true
