# deb/reverse/prerm
#!/bin/bash

systemctl stop waggle-reverse-tunnel.service || true
systemctl disable waggle-reverse-tunnel.service || true
